# .ddev/web-build/magento.cron
* * * * * cd /var/www/html && IS_DDEV_PROJECT=true php bin/magento cron:run 2>&1 | grep -v "Ran jobs by schedule" >> var/log/magento.cron.log
